# systemd/backup_hourly.timer
[Unit]
Description=BorgBackup hourly
Documentation=man:borg

[Timer]
OnCalendar=hourly
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
